# deploy/take-profit-keeper-health.service
[Unit]
Description=Take Profit keeper health API
After=network-online.target take-profit-keeper.service
Wants=network-online.target

[Service]
Type=simple
User=mongangstudio
WorkingDirectory=/home/mongangstudio/Take-Profit
Environment=KEEPER_HEALTH_PORT=8787
EnvironmentFile=-/home/mongangstudio/Take-Profit/.env
EnvironmentFile=-/home/mongangstudio/take-profit-faucet.env
ExecStart=/usr/bin/node keeper/health-server.mjs
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
